# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
redirect_stderr=true
redirect_stdout=/srv/app/logs/nginx.log

[program:php5-fpm]
command=/usr/sbin/php5-fpm
redirect_stderr=true
redirect_stdout=/srv/app/logs/php.log

[program:mysql]
command=/usr/local/bin/mysql.sh
redirect_stderr=true
redirect_stdout=/srv/app/logs/mysql.log

[program:init]
command=/init.sh
autorestart=false
redirect_stderr=true
redirect_stdout=/srv/app/logs/init.log